# automation/mac_worker/launchd/org.openpapers.prefect-worker.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.openpapers.prefect-worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/OPENPAPERS/.venv/bin/prefect</string>
    <string>worker</string>
    <string>start</string>
    <string>--pool</string>
    <string>openpapers-mac</string>
    <string>--type</string>
    <string>process</string>
    <string>--no-create-pool-if-not-found</string>
    <string>--install-policy</string>
    <string>never</string>
    <string>--no-with-healthcheck</string>
    <string>--name</string>
    <string>openpapers-mac-local</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/OPENPAPERS</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PREFECT_HOME</key>
    <string>/ABSOLUTE/PATH/TO/PREFECT/HOME</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>Umask</key>
  <integer>63</integer>
  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/LOGS/openpapers-prefect-worker.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/LOGS/openpapers-prefect-worker.stderr.log</string>
</dict>
</plist>
